# deploy/systemd/controle-financeiro-api.service
[Unit]
Description=Controle Financeiro API
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/controle-financeiro
EnvironmentFile=/var/www/controle-financeiro/.env
ExecStart=/var/www/controle-financeiro/.venv/bin/python -m uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

